# units/initrd-shutdown-journald.service.in
#  SPDX-License-Identifier: LGPL-2.1+
#
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Shutting Down Journal
DefaultDependencies=no
ConditionPathExists=/etc/initrd-release
After=initrd-switch-root.target
Before=initrd-switch-root.service

[Service]
Type=oneshot
StandardOutput=null
ExecStart=@rootbindir@/systemctl stop systemd-journald.service \
                                      systemd-journald.socket \
                                      systemd-journald-dev-log.socket
